# deploy/alpine-mirror/synchroniser/crontab
# every 5 days
* * */5 * * sh /app/mirror.sh
