# deploy/alpine-mirror/synchroniser/crontab
# crontab
*/5 * * * * root /app/mirror.sh
